# systemd/test-monitor.service
[Unit]
Description=Monitor the 'test' process and report status
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/test_monitor.sh